# plex-obsidian-sync.service.example
[Unit]
Description=Sincronizacao Plex <-> Obsidian
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/lucas/Projetos/plex-obsidian-sync
EnvironmentFile=/home/lucas/.config/plex-obsidian-sync.env
ExecStart=/home/lucas/Projetos/plex-obsidian-sync/.venv-plex-sync/bin/python /home/lucas/Projetos/plex-obsidian-sync/sync.py --loop
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
